# bootstrap/launchd/com.ma-hub.ensure-latest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ma-hub.ensure-latest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@MA_HUB_ROOT@/bootstrap/ensure-latest.sh</string>
    <string>--quiet</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>@HOUR@</integer>
    <key>Minute</key>
    <integer>@MINUTE@</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/ma-hub/ensure-latest.out.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/ma-hub/ensure-latest.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:@HOME@/bin</string>
    <key>HOME</key>
    <string>@HOME@</string>
  </dict>
</dict>
</plist>
